# services/openvpn/run
#!/bin/bash
/config/scripts/openvpnas --nodaemon --pidfile=/openvpn/pid/openvpn.pid --logfile=/config/log/openvpn.log

