# root/etc/services.d/openvpn/run
#!/usr/bin/with-contenv bash

/config/scripts/openvpnas --nodaemon --umask=0077 --pidfile=/openvpn/pid/openvpn.pid --logfile=/config/log/openvpn.log
